# deploy/systemd/personal-blog.service
[Unit]
Description=personal_blog
After=network.target

[Service]
Type=simple
WorkingDirectory=__BACKEND_DIR__
EnvironmentFile=-__BACKEND_DIR__/.env
Environment=PORT=__APP_PORT__
ExecStart=__BACKEND_DIR__/dist/personal_blog
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
